# systemd/targets/openchami.target
[Unit]
Description=Group of OpenCHAMI-related services
Requires=boot-service.service smd.service tokensmith.service metadata-service.service coresmd-coredhcp.service step-ca.service haproxy.service
Wants=coresmd-coredns.service
After=smd-init.service postgres.service coresmd-coredns.service
